# site-modules/core/files/systemd/user/opencode.service
[Unit]
Description=OpenCode local server
After=network.target

[Service]
Type=simple
Environment=TICKETS_DIR=%h/.tk
ExecStart=%h/.local/share/mise/shims/opencode web --hostname 0.0.0.0 --port 4096
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
